# pushapi/conf/supervisord.conf
[unix_http_server]
file=/letv/run/push/pushapi/supervisor.sock

[supervisord]
logfile=/letv/log/push/pushapi/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=10
loglevel=info
pidfile=/letv/run/push/pushapi/supervisord.pid
nodaemon=false
minfds=10240                ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
pushapiurl=unix:///letv/run/push/pushapi/supervisor.sock

[program:pushapi]
command=/letv/push/pushapi/pushapi -c /letv/push/pushapi/conf/conf.json
user=work
log_stdout=true
log_stderr=true
logfile=/letv/log/push/pushapi/pushapi.console.log
autorestart=true
