# pushd/conf/supervisord.conf
[unix_http_server]
file=/letv/run/push/pushd/supervisor.sock

[supervisord]
logfile=/letv/log/push/pushd/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=10
loglevel=info
pidfile=/letv/run/push/pushd/supervisord.pid
nodaemon=false
minfds=10240                ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
pushdurl=unix:///letv/run/push/pushd/supervisor.sock

[program:pushd]
command=/letv/push/pushd/pushd -c /letv/push/pushd/conf/conf.json
user=work
log_stdout=true
log_stderr=true
logfile=/letv/log/push/pushd/pushd.console.log

