# notifyapi/conf/supervisord.conf
[unix_http_server]
file=/letv/run/push/notifyapi/supervisor.sock

[supervisord]
logfile=/letv/logs/push/notifyapi/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=10
loglevel=info
pidfile=/letv/run/push/notifyapi/supervisord.pid
nodaemon=false
minfds=10240                ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
notifyapiurl=unix:///letv/run/push/notifyapi/supervisor.sock

[program:notifyapi]
command=/letv/push/notifyapi/notifyapi -c /letv/push/notifyapi/conf/conf.json
;user=work
log_stdout=true
log_stderr=true
logfile=/letv/logs/push/notifyapi/notifyapi.console.log
autorestart=true
